# mediawiki/supervisor.conf
[program:php5-fpm]
command=/usr/sbin/php5-fpm --fpm-config /etc/php5/fpm/php-fpm.conf
environment=
redirect_stderr=true
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
redirect_stderr=true
